# bin/friend_link.cron
#! /bin/sh
cd $1
$2 -f ./symfony openpne:friend-link $3
